# overlay/etc/runit/sv/hw-detect/finish
#!/bin/sh
# Remove the service symlink after successful run so it doesn't repeat
[ "$1" = "0" ] && rm -f /etc/runit/runsvdir/default/hw-detect
