# overlay/etc/runit/sv/hw-detect/run
#!/bin/sh
# runit service: NegativeOS hardware detection (runs once on first boot)
exec /usr/share/negativeos/scripts/hw-detect.sh
